# server_scripts/supervisor/heat_controller.conf
[program:heat_controller]
command = /home/heat/env/bin/python /home/heat/heat_controller.py
user = heat
stdout_logfile = /home/heat/heat_controller.log